# systemd/ha-alarm-raspy2mqtt.service
[Unit]
Description=HomeAssistant Alarm bridge for wired sensors/actuators
After=network-online.target
Wants=network-online.target

StartLimitIntervalSec=500
StartLimitBurst=5

[Service]
Type=simple
ExecStart=/root/ha-alarm-raspy2mqtt-venv/bin/python3 /root/ha-alarm-raspy2mqtt/raspy2mqtt/raspy2mqtt.py -c /etc/ha-alarm-raspy2mqtt.yaml
Environment="GPIOZERO_PIN_FACTORY=rpigpio"
Restart=on-failure
RestartSec=5s
 
[Install]
WantedBy=multi-user.target
